# overlays/cvmfs/rootfs/etc/profile.d/cvmfs.sh
#!/bin/bash

RSNT_CLUSTER=eureka
source /cvmfs/soft.computecanada.ca/config/profile/bash.sh
